# conf/supervisor.d/default/nginx.conf
[program:nginx]
command=%(ENV_HOME)s/opt/nginx/1.27.3/sbin/nginx -c %(here)s/../../nginx.d/default/nginx.conf -p %(here)s/../../../var
stdout_logfile=%(here)s/../../../var/log/%(program_name)s_stdout.log
stderr_logfile=%(here)s/../../../var/log/%(program_name)s_stderr.log
redirect_stderr=true
autostart=false
autorestart=true
