# conf/supervisor.d/default/ui.conf
[program:ui]
# debug用
# command=waitress-serve --port="8080" "ui.app":"app"

# nginx用
command=waitress-serve --unix-socket=var/run/sample_%(program_name)s.sock --asyncore-use-poll --unix-socket-perms=600 "ui.app":"app"
environment=PYTHONUNBUFFERED=1

stdout_logfile=%(here)s/../../../var/log/%(program_name)s_stdout.log
stderr_logfile=%(here)s/../../../var/log/%(program_name)s_stderr.log
redirect_stderr=true
autostart=false
autorestart=true
pcre-jit=on
stopsignal=QUIT
