# conf/supervisor.d/default/ui.conf
[program:ui]
command=%(ENV_HOME)s/opt/uwsgi/bin/uwsig --ini %(here)s/../../uwsgi.d/default/%(program_name)s.conf
stdout_logfile=%(here)s/../../../var/log/%(program_name)s_stdout.log
stderr_logfile=%(here)s/../../../var/log/%(program_name)s_stderr.log
redirect_stderr=true
autostart=false
autorestart=true
pcre-jit=on
stopsignal=QUIT
